# microservices/systemd/d31337m3-orchestrator.service
[Unit]
Description=d31337m3 Orchestrator Microservice
After=network-online.target d31337m3-auditor.service d31337m3-client-index.service d31337m3-data-handling.service d31337m3-payments.service d31337m3-watchdog.service
Wants=network-online.target

[Service]
Type=simple
User=D31337m3
WorkingDirectory=/home/D31337m3/Orm_d31337m3/microservices/orchestrator
Environment=PYTHONPATH=/home/D31337m3/Orm_d31337m3/microservices
Environment=SERVICE_PORT=8006
Environment=INFISICAL_SITE_URL=https://app.infisical.com
Environment=INFISICAL_ENVIRONMENT=prod
EnvironmentFile=-/home/D31337m3/Orm_d31337m3/.env.infisical.runtime
ExecStart=/home/D31337m3/Orm_d31337m3/microservices/.venv/bin/uvicorn service.main:app --host 127.0.0.1 --port 8006
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
